# build/s6/s6-rc.d/data-mount-check/up
/bin/sh -c "exec /etc/s6-overlay/scripts/require-data-mount.sh"
