# build/s6/s6-rc.d/labtether/run
#!/bin/sh
PG_PASS=$(cat /data/config/postgres-password)
export DATABASE_URL="postgres://${POSTGRES_USER}:${PG_PASS}@127.0.0.1:5432/${POSTGRES_DB}?sslmode=disable"
export LABTETHER_ENV=production
export LABTETHER_MANAGED_POSTGRES=true
export LABTETHER_DATA_DIR=/data
export LABTETHER_AGENT_DIR=/opt/labtether/agents
export LABTETHER_AGENT_CACHE_DIR=/data/agents
exec s6-setuidgid labtether /usr/local/bin/labtether
