# build/s6/s6-rc.d/postgres-init/up
/bin/sh -c "exec /etc/s6-overlay/scripts/postgres-init.sh"
