# build/s6/s6-rc.d/postgres/run
#!/bin/sh
mkdir -p /data/postgres /run/postgresql
chown -R postgres:postgres /data/postgres /run/postgresql
if [ ! -f /data/postgres/PG_VERSION ]; then
  /command/s6-setuidgid postgres initdb -D /data/postgres --auth-local=trust --auth-host=scram-sha-256
fi
exec /command/s6-setuidgid postgres postgres -D /data/postgres -k /run/postgresql -c listen_addresses=127.0.0.1 -c port=5432 -c unix_socket_directories=/run/postgresql
