# build/s6/s6-rc.d/guacd/run
#!/bin/sh
exec s6-setuidgid guacd /opt/guacamole/sbin/guacd -f -b 127.0.0.1 -l 4822
