# deploy/agents/launchd/com.labtether.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.labtether.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/labtether-agent</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>AGENT_NAME</key>
    <string>labtether-agent</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/var/log/labtether-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/labtether-agent.log</string>
</dict>
</plist>
